# alpine/app/lunavpn/lunavpn-lvfucs
#!/sbin/openrc-run
name="lunavpn-lvfucs"
description="LunaVPN LVFUCS Service"
command="/usr/local/bin/lvfucs"
command_args="-d /app/lunavpn/data.json"
pidfile="/var/run/lunavpn-lvfucs.pid"
start_stop_daemon_args="--background --make-pidfile"
depends="docker-compose-wg"

depend() {
  need docker-compose-wg
}

stop_post() {
  # Clean up after the program exits
  rm -f "${pidfile}"
}
